# debian/eos-kalite-systemd-helper.socket
[Unit]
Description=KA Lite Server

[Socket]
ListenStream=8008
BindToDevice=lo

[Install]
WantedBy=default.target
