# debian/eos-kalite-systemd-helper.socket
[Unit]
Description=KALite HTTPd

[Socket]
ListenStream=8008

[Install]
WantedBy=default.target
