# debian/kalite-systemd-helper.service
[Unit]
Description=KALite HTTPd

[Service]
ExecStart=/usr/bin/kalite-flatpak-wrapper
User=kalite

PrivateTmp=yes
